# vps/config/frp.conf
[program:frp]
command = /home/frp/frps -c /home/frp/frps.ini
autostart = true

